# scripts/vhdx_mount.service
[Unit]
Description=Mount VHDX files for directories
Before=network-pre.target

[Service]
ExecStart=/usr/local/sbin/vhdx_mount.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=local-fs.target
